# root/etc/services.d/quassel/run
#!/usr/bin/with-contenv bash

exec \
	s6-setuidgid abc quasselcore \
	-c /config -p 4242 --require-ssl
